# workshop/supervisor/metrics.conf
[program:metrics]
process_name=metrics
autostart=false
command=/home/eduk8s/bin/start-metrics
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
